# .config/systemd/user/mpd.service.d/override.conf
[Unit]
Wants=mpd-mpris.service mpdas.service
