# .config/systemd/user/mpdas.service.d/override.conf
[Unit]
PartOf=mpd.service
BindsTo=mpd.service
